# assets/alpine/myst-service
#!/sbin/openrc-run

if [ -f /tmp/.env.myst ]; then
  source /tmp/.env.myst
fi

name="Mysterium Network Node"
command="myst"
command_args="--keystore.lightweight --local-service-discovery=false --launcher.ver=$LAUNCHER service --agreed-terms-and-conditions"
command_background=true
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
  need net
}

start_pre() {
  if ! [ -f /tmp/.env.myst ]; then
    echo "Error: /tmp/.env.myst is absent"
    exit 1
  fi

  l=`find /root/.mysterium/keystore -type f 2>/dev/null | wc -l`
  if [ $l -ge 2 ]
  then
    logger Keystore found..
  else
    echo "Error: keystore not found" 
    exit 1
  fi
}
