# assets/alpine/myst-service
#!/sbin/openrc-run

name="Mysterium Network Node"
command="vm-myst-agent"
command_args=""
command_background=true
pidfile="/run/${RC_SVCNAME}.pid"
#output_log="/var/log/vm-myst-agent.log"
#error_log="/var/log/vm-myst-agent.err"

depend() {
  need net
}
