# assets/alpine/myst-service
#!/sbin/openrc-run

name="Mysterium Network Node"
command="vm-agent"
command_args=""
pidfile="/run/vm-myst-agent.pid"

depend() {
  need net
}
